# systemd/steinsgate-wallpaper.timer
[Unit]
Description=Rotate Steins;Gate wallpaper every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
