# systemd/steinsgate-wallpaper.service
[Unit]
Description=Rotate Steins;Gate wallpaper every hour

[Service]
Type=oneshot
ExecStart=%h/rotate-wallpaper.sh
